# conf/storm/supervisord-worker.conf
[program:storm-supervisor]
command=storm supervisor
directory=/usr/lib/storm
autorestart=true
user=storm

[program:storm-logviewer]
command=storm logviewer
autorestart=true
user=storm
